# fff_api/rootfs/etc/services.d/fff-api/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# JPD Hass.io Add-ons: FFF-API Server
# Take down the S6 supervision tree when fff-api fails
# ==============================================================================


s6-svscanctl -t /var/run/s6/services
